# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
echo 60 > /sys/class/hwmon/hwmon1/pwm1

exit 0
